# tg-bot/tg-bot.service
# systemd service file for deployment
# /etc/systemd/system/tg-bot.service

[Unit]
Description=Telegram Bot
After=network.target

[Service]
User=your_user
WorkingDirectory=/home/your_user/tg-bot
ExecStart=/home/your_user/tg-bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target